# ubuntu18.04-0/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:php]
command=/usr/sbin/service php7.4-fpm start

